# backend/crontab
# Daily Horse Racing Crawler Crontab
# Runs at 1:00 AM every day to crawl previous day's races
#
# Format: minute hour day month day-of-week command
# All times in UTC inside container (adjust if needed)

# Run daily at 1:00 AM (adjust timezone as needed)
0 1 * * * cd /app && /usr/local/bin/python3 /app/backend/daily_crawl.py >> /var/log/horse-racing-crawler.log 2>&1

# Run Upcoming Entries Crawler at 5:05 AM (12:05 AM EST) and 10:05 AM (5:05 AM EST)
5 5,10 * * * cd /app && /usr/local/bin/python3 /app/backend/crawl_entries.py >> /var/log/horse-racing-crawler.log 2>&1

# Run Live Results Crawler every 30 mins from 11 AM to 11 PM (for TODAY's races)
*/30 11-23 * * * cd /app && /usr/local/bin/python3 /app/backend/crawl_equibase.py $(date +\%Y-\%m-\%d) >> /var/log/horse-racing-crawler.log 2>&1

# Keep container alive - log heartbeat
*/30 * * * * echo "[$(date)] Scheduler heartbeat" >> /var/log/horse-racing-crawler.log

# Empty line required at end of crontab

